# deploy/appSetup/hollowapp.conf
[program:hollowapp]
command=/home/ubuntu/hollowapp/venv/bin/gunicorn -b localhost:8000 -w 4 hollowapp:app
directory=/home/ubuntu/hollowapp
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true